# wgen/wgen.cron
MAILTO=""
55 2,8,14,20 * * * ec2-user /home/ec2-user/projects/aqm/wgen/bin/python /home/ec2-user/projects/aqm/wgen/wgen.py -w
0 0 * * *          ec2-user /home/ec2-user/projects/aqm/wgen/bin/python /home/ec2-user/projects/aqm/wgen/wgen.py -i
